# systemd/epaper26.timer
[Unit]
Description=Timer to run epaper26 service regularly

[Timer]
OnCalendar=*:*:0
Persistent=true

[Install]
WantedBy=timers.target
